# docker/symfony/root/etc/services.d/api/run
#!/usr/bin/with-contenv bash
cd /opt/backend
exec symfony server:start
