# docker/fastapi/root/etc/services.d/api/run
#!/usr/bin/env bash

cd /opt/backend
. ./venv/bin/activate
PYTHONPATH=. HOST=127.0.0.1 PORT=8000 python3 -m pyfoods
